# moo/moo.service
[Unit]
Description=MOO server

[Service]
Type=simple
User=moo
Group=moo
ExecStart=/home/moo/bin/run-moo

[Install]
WantedBy=multi-user.target
